# personal/crontab
*/2 * * * * ~/ansel/personal/ansel-entry.sh
